# xserver/config/openclaw-gateway-eleven.service
[Unit]
Description=OpenClaw Gateway Eleven (v2026.3.2)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/root/.nvm/versions/node/v24.13.1/bin/node /root/.nvm/versions/node/v24.13.1/bin/../lib/node_modules/openclaw/dist/index.js gateway --port 18791 --allow-unconfigured
Restart=always
RestartSec=5
KillMode=process
Environment=HOME=/root
Environment=OPENCLAW_STATE_DIR=/root/.openclaw-eleven
Environment=OPENCLAW_PROFILE=eleven
Environment=OPENCLAW_GATEWAY_PORT=18791
Environment="OPENCLAW_SYSTEMD_UNIT=openclaw-gateway-eleven.service"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.3.2
Environment=ZAI_API_KEY=${ZAI_API_KEY}
Environment=DISCORD_BOT_TOKEN=${DISCORD_BOT_TOKEN}

[Install]
WantedBy=default.target
