# xserver/config/openclaw-gateway-alpha.service
# OpenClaw Gateway Alpha systemd service template
# Discord bot instance (OpenRex-alpha, 常識役) on port 18791
# Variables to be replaced during restore:
#   {{NODE_PATH}} -> /root/.nvm/versions/node/v24.13.1/bin/node
#   {{OPENCLAW_PATH}} -> /root/.nvm/versions/node/v24.13.1/lib/node_modules/openclaw/dist/index.js

[Unit]
Description=OpenClaw Gateway Alpha (v2026.3.2)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart={{NODE_PATH}} {{OPENCLAW_PATH}} gateway --port 18791 --allow-unconfigured
Restart=always
RestartSec=5
KillMode=process
Environment=HOME=/root
Environment=OPENCLAW_STATE_DIR=/root/.openclaw-alpha
Environment=OPENCLAW_PROFILE=alpha
Environment=OPENCLAW_GATEWAY_PORT=18791
Environment="OPENCLAW_SYSTEMD_UNIT=openclaw-gateway-alpha.service"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.3.2
Environment=ZAI_API_KEY={{ZAI_API_KEY}}
Environment=DISCORD_BOT_TOKEN={{DISCORD_BOT_TOKEN}}

[Install]
WantedBy=default.target
