# xserver/config/openclaw-gateway.service
# OpenClaw Gateway systemd service template
# Variables to be replaced during restore:
#   {{NODE_PATH}} -> /root/.nvm/versions/node/v24.13.1/bin/node
#   {{OPENCLAW_PATH}} -> /root/.nvm/versions/node/v24.13.1/lib/node_modules/openclaw/dist/index.js

[Unit]
Description=OpenClaw Gateway (v2026.3.2)
After=network-online.target
Wants=network-online.target

[Service]
# Auto-apply LINE restart loop patch before starting (survives npm update -g openclaw)
ExecStartPre=/root/openclaw-scripts/auto-line-patch.sh
ExecStart={{NODE_PATH}} {{OPENCLAW_PATH}} gateway --port 18789 --allow-unconfigured
Restart=always
RestartSec=5
KillMode=process
Environment=HOME=/root
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment="OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.3.2
Environment=ZAI_API_KEY={{ZAI_API_KEY}}

[Install]
WantedBy=default.target
